# ansible/roles/aws-bc-tools/files/ec2-local-storage/bc-ec2-local-storage.service
[Unit]
Description=EC2 Route 53 Boot Script
After=cloud-final.service

[Service]
User=bradley
Group=bradley
ExecStart=/opt/bc/ec2-r53-boot.py
Type=oneshot

[Install]
WantedBy=multi-user.target
